# templates/build-agent.erb
#!/bin/sh -
# File is managed by Puppet

### BEGIN INIT INFO
# Provides:          build-agent
# Required-Start:    $local_fs $syslog
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Build agent management script
### END INIT INFO

. /lib/lsb/init-functions

NAME=build-agent
DAEMON=<%= @agent_dir %>/bin/agent.sh
PIDFILE=<%= @agent_dir %>/logs/buildAgent.pid
USER=<%= @agent_user %>

# If the daemon is not there, then exit.
test -x $DAEMON || exit 5

case $1 in
start)
  if [ -e $PIDFILE ]; then
    status_of_proc -p $PIDFILE $DAEMON "$NAME process" >/dev/null && status="0" || status="$?"
    if [ $status = "0" ]; then
      log_daemon_msg "$NAME process is already running"
      exit
    fi
  fi

  log_daemon_msg "Starting the process" "$NAME"
  su --login -c "$DAEMON start" $USER || log_end_msg 1
;;
stop)
  if [ -e $PIDFILE ]; then
    status_of_proc -p $PIDFILE $DAEMON "Stoppping the $NAME process" >/dev/null && status="0" || status="$?"
    if [ $status = "0" ]; then
      su --login -c "$DAEMON stop" $USER || log_end_msg 1
    else
      log_daemon_msg "$NAME process is not running"
      log_end_msg 0
    fi
  fi
;;
restart)
  $0 stop && sleep 2 && $0 start
  ;;
status)
  if [ -e $PIDFILE ]; then
    status_of_proc -p $PIDFILE $DAEMON "$NAME process" && exit 0 || exit $?
  else
    log_daemon_msg "$NAME Process is not running"
    log_end_msg 0
  fi
;;
*)
  echo "Usage: $0 {start|stop|restart|status}"
  exit 2
esac
