# templates/build-agent.erb
#!/bin/sh -
# File is managed by Puppet

### BEGIN INIT INFO
# Provides:          build-agent
# Required-Start:    $local_fs $syslog
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Build agent management script
### END INIT INFO

export TEAMCITY_AGENT_MEM_OPTS="<%= @teamcity_agent_mem_opts %>"

su --login -c "<%= @destination_dir %>/<%= @agent_dir %>/bin/agent.sh '$@'" <%= @agent_user %>
