# templates/build-agent.erb
#!/bin/sh -
# File is managed by Puppet

### BEGIN INIT INFO
# Provides:          build-agent
# Required-Start:    $local_fs $syslog
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Build agent management script
### END INIT INFO

su --login -c "/var/tainted/build-agent/bin/agent.sh '$@'" teamcity
